# umpleonline/docker_config/supervisord.conf
[supervisord]
# This was originally nginx, but having trouble starting Docker in Docker
user=root
loglevel=debug

[program:nginx]
command=nginx -g 'daemon off;'

[program:php-fpm]
command=php-fpm7 --nodaemonize

[program:codeExecNodeInDocker]
user=root
command=bash /var/UmpleCodeExecution/setup.sh

[program:lsp-proxy]
command=node /var/www/scripts/lsp-proxy/server.js
directory=/var/www/scripts/lsp-proxy
user=nginx
autostart=true
autorestart=true
stderr_logfile=/var/log/lsp-proxy.err.log
stdout_logfile=/var/log/lsp-proxy.out.log
environment=UMP_BASE_DIR="/var/www/ump",UMPLE_LSP_WS_URL="/lsp",LSP_AUTH_SECRET="%(ENV_LSP_AUTH_SECRET)s",LSP_DEV_AUTH="1",UMPLESYNC_JAR_PATH="/var/www/scripts/umplesync.jar"

