# umpleonline/docker_config/supervisord.conf
[supervisord]
user=nginx
loglevel=debug

[program:nginx]
command=nginx -g 'daemon off;'

[program:php-fpm]
command=php-fpm7 --nodaemonize

